# agent/dashboard.service
[Unit]
Description=PSBA AICC Web Dashboard
After=network.target aiagent.service saima.service

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=/opt/aiagent
ExecStart=/opt/aiagent/venv/bin/python3 /opt/aiagent/dashboard.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
