# agent/sana.service
[Unit]
Description=AI Sales Agent — Sana (Urdu Sales, ext 8500)
After=network.target asterisk.service zara.service

[Service]
Type=simple
User=asterisk
WorkingDirectory=/opt/aiagent
ExecStart=/opt/aiagent/venv/bin/python3 /opt/aiagent/sana.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sana

[Install]
WantedBy=multi-user.target
